# infra/systemd/rafael-lift.service
[Unit]
Description=rafael lift
After=network-online.target

[Service]
Type=simple

Environment=DATABASE_URL=sqlite:%h/rafael/data/lift/lift.db?mode=rwc
Environment=LIFT_STATIC_DIR=%h/rafael/services/lift/frontend/dist
Environment=LIFT_BIND=127.0.0.1:3033

WorkingDirectory=%h/rafael
ExecStartPre=mkdir -p %h/rafael/data/lift
ExecStart=%h/rafael/target/release/lift

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
