# infra/systemd/rafael-chat.service
# Configuration is done via environment variables.
# Override defaults by creating a drop-in file:
#
#   ~/.config/systemd/user/rafael-chat.service.d/override.conf
#
# With contents like:
#
#   [Service]
#   Environment=RAFAEL_CHAT_BIND=127.0.0.1:3031
#   Environment=RAFAEL_CHAT_MODEL_BASE_URL=http://rafael:8080/v1
#   Environment=RAFAEL_CHAT_MODEL_NAME=gemma-everyday
#   Environment=RAFAEL_CHAT_MODEL_LIST_TIMEOUT_SECONDS=5
#   Environment=RAFAEL_CHAT_AUTH_TOKEN_DAYS=30
#   Environment=RAFAEL_CHAT_WEB_SEARCH_PROVIDER=searxng
#   Environment=RAFAEL_CHAT_SEARXNG_BASE_URL=http://127.0.0.1:8888/
#   Environment=RAFAEL_CHAT_WEB_SEARCH_FETCH_RESULTS=3
#   Environment=RAFAEL_CHAT_WEB_SEARCH_FETCH_MAX_BYTES=16384
#
# Or for Brave Search:
#
#   [Service]
#   Environment=RAFAEL_CHAT_WEB_SEARCH_PROVIDER=brave
#   Environment=RAFAEL_CHAT_BRAVE_SEARCH_API_KEY=...
#
# Or run: systemctl --user edit rafael-chat
#
[Unit]
Description=rafael chat
Wants=llama-swap.service
After=network-online.target llama-swap.service

[Service]
Type=simple

Environment=RAFAEL_CHAT_BIND=0.0.0.0:3031
Environment=RAFAEL_CHAT_DATA_DIR=%h/rafael/data/chat
Environment=RAFAEL_CHAT_WEB_DIST=%h/rafael/services/chat/web/dist
Environment=RAFAEL_CHAT_MODEL_BASE_URL=http://rafael:8080/v1
Environment=RAFAEL_CHAT_MODEL_NAME=gemma-everyday
Environment=RAFAEL_CHAT_MODEL_TIMEOUT_SECONDS=300
Environment=RAFAEL_CHAT_MODEL_LIST_TIMEOUT_SECONDS=5
Environment=RAFAEL_CHAT_AUTH_TOKEN_DAYS=30

WorkingDirectory=%h/rafael
ExecStart=%h/rafael/target/release/chat serve

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
