# infra/systemd/rafael-coding.service
[Unit]
Description=Run the coding service
After=network.target

[Service]
Type=simple
User=cachebag
ExecStart=/home/cachebag/rafael/target/release/coding serve
EnvironmentFile=/home/cachebag/.config/rafael/coding.env
Restart=on-failure

[Install]
WantedBy=default.target